# docker/php/production/configs/supervisord.conf
[supervisord]
nodaemon=true
chmod=0770
user=www-data

[program:php-fpm]
command = /usr/local/sbin/php-fpm
autostart=true
autorestart=unexpected
exitcodes=0

[supervisorctl]
